# device/rom-sorter.timer
[Unit]
Description=DeckDock — Run ROM Sorter every 30 seconds

[Timer]
OnBootSec=60
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
