# device/save-backup-watcher.service
[Unit]
Description=DeckDock — Watch for suspend and backup saves
After=default.target

[Service]
Type=simple
ExecStart=/bin/bash %h/Emulation/tools/sleep-watcher.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
